# little-genesis.service
[Unit]
Description=Little Genesis Start Script

[Service]
ExecStart=/projects/little-genesis-python/lg.py

[Install]
WantedBy=multi-user.target 
